# udev/rules.d/10-monitor.rules
ACTION=="change", SUBSYSTEM=="drm", RUN+="/etc/udev/scripts/monitor-hotplug.sh %p"
