# udev/rules.d/10-monitor.rules
ACTION=="change", SUBSYSTEM=="drm", RUN+="/usr/local/bin/monitor-hotplug.sh %p"
